# db/crontab
0 0 * * * /backups/backup.sh >> /var/log/pg_backup.log 2>&1
